# deploy/systemd/docutrust-blockchain.service
[Unit]
Description=DocuTrust Blockchain Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/docutrust/current/blockchain-service
Environment=NODE_ENV=production
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

